# _common-files/etc/systemd/system/genieacs-fs.service
[Unit]
Description=GenieACS FS
After=network.target mongod.service
Requires=mongod.service

[Service]
User=genieacs
EnvironmentFile=/opt/genieacs/genieacs.env
ExecStart=/usr/bin/genieacs-fs
Restart=on-failure

[Install]
WantedBy=multi-user.target
